# vcontrold/rootfs/etc/services.d/vclient_pub/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

## Fetch MQTT settings
. /etc/services.d/get_mqtt_settings.sh

## Fetch Vcontrold settings
. /etc/services.d/get_vcontrold_settings.sh

## Run your program
bashio::log.info "Starting vclient_pub main loop"
while sleep $REFRESH_RATE; do
    vclient -h $VCONTROL_HOST -p $VCONTROL_PORT -f /etc/vcontrold/1_mqtt_commands.txt -t /etc/vcontrold/2_mqtt.tmpl -x /etc/vcontrold/3_mqtt_pub.sh && exitcode=0 || exitcode=$?
    bashio::log.info "vclient_pub run finished with exit code $exitcode. Looping vclient. Sleeping $REFRESH_RATE seconds."
done
