# vcontrold/rootfs/etc/services.d/vclient_pub/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the vclient_pub service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Fetch MQTT settings
source /etc/services.d/get_mqtt_settings.sh

# Fetch Vcontrold settings
source /etc/services.d/get_vcontrold_settings.sh

# Ensure required variables are set
if [[ -z "$REFRESH_RATE" || -z "$VCONTROL_HOST" || -z "$VCONTROL_PORT" ]]; then
    bashio::log.error "Missing essential environment variables. Check MQTT and Vcontrold settings."
    exit 1
fi

# Run your program
while true; do
    bashio::log.info "Executing vclient_pub command..."
    if vclient -h "$VCONTROL_HOST" -p "$VCONTROL_PORT" -f /etc/vcontrold/1_mqtt_commands.txt \
               -t /etc/vcontrold/2_mqtt.tmpl -x /etc/vcontrold/3_mqtt_pub.sh; then
        exitcode=0
        bashio::log.info "vclient_pub executed successfully."
    else
        exitcode=$?
        bashio::log.error "vclient_pub run failed with exit code $exitcode."
    fi
    bashio::log.info "Looping vclient_pub. Sleeping for $REFRESH_RATE seconds."
    sleep "$REFRESH_RATE"
done
