# vcontrold/rootfs/etc/services.d/vclient_sub/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

declare refresh_rate
refresh_rate=$(bashio::config 'refresh')

## Fetch MQTT credentials
. /etc/services.d/get_mqtt_credentials.sh


bashio::log.info "Setting vcontrold host and port ..."
export VCONTROL_HOST=$(bashio::config 'vcontrol_host')
export VCONTROL_PORT=$(bashio::config 'vcontrol_port')

## Run your program
while true ; do
    mosquitto_sub -v -h $MQTT_HOST -p $MQTT_PORT -u $MQTT_USER -P $MQTT_PASSWORD -t "$MQTT_TOPIC/#" | while read -r payload
    do
        # Here is the callback to execute whenever you receive a message:
        topic=$(echo $payload | cut -d' ' -f1)
        value=$(echo $payload | cut -d' ' -f2)
        if [[ $topic =~ "/set" ]]; then
            command=$(echo $topic | cut -d'/' -f2)
            bashio::log.info "Sending command: [${command}]: ${value}"
            vclient -h $VCONTROL_HOST -p $VCONTROL_PORT -o /dev/stdout -c "${command} ${value}"
        fi
    done
    sleep $refresh_rate
done
